# systemd_files/chakna-audio-speaker.service
[Unit]
Description=Chakna Speaker Service
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/chakna
ExecStart=/home/pi/.local/bin/uv run -m actuators.audio.speaker_service
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
